# pi-code/restartScript.conf
# put it in /etc/init/

start on startup
start on runlevel [2345]
stop on runlevel [016]

# respawn
exec python /path/to/your/script.py